# openclaw-infra/configs/openclaw.service
# ============================================================
# OpenClaw Gateway systemd 服务文件
#
# 安装位置：/etc/systemd/system/openclaw.service
#
# 作用：
# - 系统启动后自动拉起 OpenClaw Gateway
# - 进程崩溃后 5 秒自动重启
# - API Key 从 EnvironmentFile 安全注入
#
# 安装命令：
#   sudo cp openclaw.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now openclaw
#
# 前置条件：
#   - 已通过 npm install -g openclaw@latest 安装 OpenClaw
#   - 已创建 /opt/openclaw.env 并填入 API Key
# ============================================================

[Unit]
Description=OpenClaw Gateway
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

# API Key 等敏感信息从独立文件注入
EnvironmentFile=/opt/openclaw.env

# OpenClaw 需要 HOME 目录读取 ~/.openclaw 配置
Environment=HOME=/root

# 启动 gateway，绑定 18789 端口（默认 loopback，仅本机可访问）
ExecStart=/usr/bin/openclaw gateway --port 18789

# 无论什么原因退出，5 秒后自动重启
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
